# yocto/profiles/production/etc/systemd/system/dirtsim-server.service.d/restart.conf
# Production profile: auto-restart on failure.
[Service]
Restart=on-failure
RestartSec=5
