# yocto/meta-dirtsim/recipes-dirtsim/dirtsim/files/dirtsim-config-setup.service
[Unit]
Description=DirtSim Configuration Setup
DefaultDependencies=no
RequiresMountsFor=/data/dirtsim
After=local-fs.target data.mount
Before=dirtsim-audio.service dirtsim-os-manager.service dirtsim-server.service dirtsim-ui.service

[Service]
Type=oneshot
ExecStart=/usr/bin/dirtsim-config-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
