# yocto/meta-dirtsim/recipes-dirtsim/dirtsim/files/dirtsim-ui.service
[Unit]
Description=DirtSim Simulation UI
Wants=dirtsim-audio.service dirtsim-server.service
Requires=dirtsim-config-setup.service
After=dirtsim-config-setup.service

[Service]
Type=simple
User=dirtsim
Group=video
WorkingDirectory=/data/dirtsim

# Unbind kernel console from framebuffer so dmesg doesn't draw over LVGL.
# Pi4 has vtcon1 (framebuffer console), Pi5 only has vtcon0 (dummy).
# Needs to run as root (+ prefix) since vtconsole sysfs requires elevated privileges.
ExecStartPre=+/bin/sh -c 'test -e /sys/class/vtconsole/vtcon1/bind && echo 0 > /sys/class/vtconsole/vtcon1/bind || true'

ExecStart=/usr/bin/dirtsim-ui -b drm --connect localhost:8080
Restart=no

# Display settings are auto-detected by dirtsim-detect-display.service,
# which writes /etc/dirtsim/display.conf based on connected hardware.
# Defaults below are used if the conf file doesn't exist.
Environment=LV_LINUX_DRM_DEVICE=/dev/dri/card2
Environment=LV_DISPLAY_ROTATION=270
Environment=LV_EVDEV_DEVICE=/dev/input/by-path/platform-i2c@0-event
EnvironmentFile=-/etc/dirtsim/display.conf

# Give access to framebuffer.
SupplementaryGroups=video input

[Install]
WantedBy=multi-user.target
