# yocto/meta-dirtsim/recipes-dirtsim/dirtsim/files/dirtsim-server.service
[Unit]
Description=DirtSim Physics Simulation Server
After=avahi-daemon.service

[Service]
Type=simple
User=dirtsim
WorkingDirectory=/data/dirtsim
ExecStart=/usr/bin/dirtsim-server -p 8080
Restart=no

[Install]
WantedBy=multi-user.target
