# yocto/meta-dirtsim/recipes-dirtsim/sparkle-duck/files/sparkle-duck-server.service
[Unit]
Description=Sparkle Duck Physics Simulation Server
After=network-online.target avahi-daemon.service
Wants=network-online.target

[Service]
Type=simple
User=dirtsim
WorkingDirectory=/home/dirtsim/sparkle-duck
ExecStart=/usr/bin/sparkle-duck-server -p 8080
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
