# yocto/meta-dirtsim/recipes-dirtsim/sparkle-duck/files/sparkle-duck-detect-display.service
[Unit]
Description=Detect display type and configure sparkle-duck-ui
DefaultDependencies=no
Before=sparkle-duck-ui.service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/sparkle-duck-detect-display.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
