# yocto/meta-dirtsim/recipes-dirtsim/sparkle-duck/files/sparkle-duck-ui.service
[Unit]
Description=Sparkle Duck Simulation UI
After=sparkle-duck-server.service
Requires=sparkle-duck-server.service

[Service]
Type=simple
User=dirtsim
Group=video
WorkingDirectory=/home/dirtsim/sparkle-duck

# Unbind kernel console from framebuffer so dmesg doesn't draw over LVGL.
# Needs to run as root (+ prefix) since vtconsole sysfs requires elevated privileges.
ExecStartPre=+/bin/sh -c 'echo 0 > /sys/class/vtconsole/vtcon1/bind'

ExecStart=/usr/bin/sparkle-duck-ui -b fbdev --connect localhost:8080
Restart=always
RestartSec=5

# Display rotation for HyperPixel 4.0 (landscape orientation).
# LVGL handles rotation in software - no hardware dtparam needed.
Environment=LV_DISPLAY_ROTATION=90

# Touchscreen input device (Goodix capacitive touch on HyperPixel 4.0).
# Using stable by-path symlink instead of eventN which changes on reboot.
Environment=LV_EVDEV_DEVICE=/dev/input/by-path/platform-i2c@0-event

# Give access to framebuffer.
SupplementaryGroups=video input

[Install]
WantedBy=multi-user.target
